# dot_config/systemd/user/chezmoi-pull.timer
[Unit]
Description=Weekly chezmoi pull and apply timer
Requires=chezmoi-pull.service

[Timer]
# Run every Sunday at 2:00 AM
OnCalendar=Sun *-*-* 02:00:00
# Run on boot if missed
Persistent=true

[Install]
WantedBy=timers.target
